# addons/etc/init.d/loki
#!/sbin/openrc-run
description="Loki (daemonized) Simple IOC And Yara Scanner (fork)"
extra_commands=""
extra_started_commands="update"

workdir=/opt/Loki-daemonized
pidfile=${workdir}/loki.pid
command=/opt/Loki-daemonized/bin/python3
command_args="loki.py -s 20000 -d --noprocscan --noindicator --dontwait --intense --csv --nolog --force"

depend() {
	need net
	use dns logger netmount
}

start() {
	ebegin "Starting Loki daemon"
	cd $workdir
	source bin/activate
	$command $command_args | logger -t loki &
}

stop() {
	ebegin "Stopping Loki daemon"
	kill -TERM $(cat $pidfile) &> /dev/null &
}

restart() {
	stop
	start
}

update() {
	cd $workdir
	stop
	. bin/activate
	./loki-upgrader.py --sigsonly --nolog
	deactivate
	start
}
