# addons/etc/init.d/paranoya
#!/sbin/openrc-run
description="paranoya: Simple IOC And Yara Scanner"
extra_commands=""
extra_started_commands="update"

workdir=/opt/paranoya
pidfile=${workdir}/paranoya.pid
command=${workdir}/paranoya.py
command_args="-s 20000 -d --noprocscan --noindicator --intense --csv --nolog --force"

depend() {
	need net
	use dns logger netmount
}

start() {
	ebegin "Starting paranoya daemon"
	cd $workdir
	$command $command_args | logger -t paranoya &
}

stop() {
	ebegin "Stopping paranoya daemon"
	kill -TERM $(cat $pidfile) &> /dev/null &
}

restart() {
	stop
	start
}

update() {
	cd $workdir
	stop
	./upgrader.py --sigsonly
	start
}
